# gpu-monitor.service
[Unit]
Description=GPU Temperature Monitor
After=network.target

[Service]
Type=simple
User=$USER
WorkingDirectory=$HOME/.local/share/gpu-monitor
Environment=PATH=$HOME/.local/share/gpu-monitor/venv/bin:$HOME/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=$HOME/.local/share/gpu-monitor/venv/bin/python $HOME/.local/share/gpu-monitor/gpu_monitor.py --service
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target 